# launch_agents/com.user.telegram-bot.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.telegram-bot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/path/to/garmin-drive-sync/venv/bin/python3</string>
    <string>/path/to/garmin-drive-sync/scripts/telegram_bot.py</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/path/to/garmin-drive-sync</string>

  <key>StandardOutPath</key>
  <string>/path/to/garmin-drive-sync/logs/telegram_bot_stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/garmin-drive-sync/logs/telegram_bot_stderr.log</string>
</dict>
</plist>
